# packaging/ac.socket
[Socket]
ListenStream=/run/aul/daemons/%U/.amd-sock
DirectoryMode=0777

[Install]
WantedBy=sockets.target
