# packaging/ac.socket
[Socket]
ListenStream=/run/amd/%U
DirectoryMode=0777

[Install]
WantedBy=sockets.target
